# ltafantasybot.service
[Unit]
Description=LTA Fantasy Bot - Telegram Bot for League of Legends Fantasy
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/bots/ltafantasybot
Environment="PATH=/home/ubuntu/bots/ltafantasybot/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# Run the bot
ExecStart=/home/ubuntu/bots/ltafantasybot/.venv/bin/python /home/ubuntu/bots/ltafantasybot/bot.py

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=0

# Logging
StandardOutput=append:/home/ubuntu/bots/ltafantasybot/bot.log
StandardError=append:/home/ubuntu/bots/ltafantasybot/bot.log

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Resource limits (optional)
# LimitNOFILE=4096
# MemoryMax=512M

[Install]
WantedBy=multi-user.target
